# services.d/tensorboard/run
#!/usr/bin/execlineb -P
/usr/local/bin/tensorboard --logdir /data/logs
